# systemd/willow.service
[Unit]
Description=Willow Voice Assistant Service
Documentation=https://github.com/saim/willow
After=graphical-session.target pipewire.service pipewire-pulse.service sound-server.target
Wants=pipewire-pulse.service sound-server.target
PartOf=graphical-session.target

[Service]
Type=dbus
BusName=com.github.saim.Willow
ExecStart=/usr/bin/willow-service
TimeoutStopSec=15
Restart=on-failure
RestartSec=5

Environment=PATH=/usr/bin:/usr/local/bin
Environment=XDG_RUNTIME_DIR=/run/user/%U
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/%U/bus

# Import graphical session environment for xdg-open / browser launch
ExecStartPre=/bin/sh -c 'for var in WAYLAND_DISPLAY DISPLAY XDG_CURRENT_DESKTOP; do eval val=$(tr "\\0" "\\n" < /proc/$(pgrep -u %U -n gnome-shell 2>/dev/null || echo 0)/environ 2>/dev/null | grep "^$var=" | cut -d= -f2-); [ -n "$val" ] && systemctl --user set-environment $var="$val"; done; exit 0'

# Security - Relaxed for GUI app launching
PrivateTmp=false
NoNewPrivileges=false

# Auto-create ~/.config/willow on first start (ReadWritePaths fails if missing)
ConfigurationDirectory=willow

[Install]
WantedBy=graphical-session.target
